# deploy/sensitive-word.service
[Unit]
Description=Sensitive Word Risk Review Tool
After=network.target

[Service]
User=root
WorkingDirectory=/opt/sensitive-word
EnvironmentFile=/opt/sensitive-word/.env
ExecStart=/opt/sensitive-word/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8002 \
    --workers 1 \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
